# rootfs/etc/service/php-fpm/run
#!/bin/sh

/usr/local/sbin/php-fpm --fpm-config /usr/local/etc/php-fpm.d/www.conf --allow-to-run-as-root --nodaemonize