# rootfs/etc/service/php-fpm/run
#!/bin/sh

/usr/local/sbin/php-fpm -F --nodaemonize