# 7.3/rootfs/etc/service/nginx/finish
#!/bin/sh

exec kill `cat /run/nginx/nginx.pid`
